# etc/init.d/secretservice
#!/sbin/openrc-run

name="secretservice"
description="Secret Sharing Service"
command="/usr/bin/secretservice"
command_background="yes"
pidfile="/run/secretservice.pid"
user="secretservice"
group="secretservice"
output_log="/var/log/secretservice.log"
error_log="/var/log/secretservice.log"

directory="${SS_PATH}"
umask="0077"

depend() {
    need net
    use dns
}

start_pre() {
    checkpath --directory --owner "${user}:${group}" --mode "0750" "${SS_PATH}"

    # Set capabilities for privileged ports if not in proxy mode
    if [ "${SS_PROXY}" = "false" ]; then
        setcap cap_net_bind_service=+ep /usr/bin/secretservice 2>/dev/null || \
            ewarn "Could not set capabilities - may need to run as root or setcap manually"
    fi
}

stop_timeout=65

stop() {
    kill -TERM $(cat "${pidfile}") 2>/dev/null
    # Wait for process to actually exit (graceful shutdown takes up to 1 minute)
    local pid
    pid=$(cat "${pidfile}" 2>/dev/null)
    if [ -n "$pid" ]; then
        while kill -0 "$pid" 2>/dev/null; do
            sleep 1
        done
    fi
}
